# products/Private-AI-Lab/apps/jvt-inbound-voice-agent/com.jvt.inbound-voice-agent.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jvt.inbound-voice-agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/c.s.d.v.r.s./Developer/Control-Host/JVT-Technologies/products/Private-AI-Lab/apps/jvt-inbound-voice-agent/tools/run_voice_agent.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/c.s.d.v.r.s./Developer/Control-Host/JVT-Technologies/products/Private-AI-Lab/apps/jvt-inbound-voice-agent/data/voice-agent.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/c.s.d.v.r.s./Developer/Control-Host/JVT-Technologies/products/Private-AI-Lab/apps/jvt-inbound-voice-agent/data/voice-agent.stderr.log</string>
</dict>
</plist>

