# ops/control-panel/com.jvt.control-panel.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jvt.control-panel</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/c.s.d.v.r.s./Developer/Control-Host/JVT-Technologies/ops/control-panel/run_control_panel.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/c.s.d.v.r.s./Developer/Control-Host/JVT-Technologies</string>
    <key>StandardOutPath</key>
    <string>/Users/c.s.d.v.r.s./Developer/Control-Host/JVT-Technologies/ops/control-panel/control-panel.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/c.s.d.v.r.s./Developer/Control-Host/JVT-Technologies/ops/control-panel/control-panel.stderr.log</string>
  </dict>
</plist>
